# extras/social-relay.target
[Unit]
Description=social-relay
Wants=redis.service
After=redis.service

[Install]
WantedBy=multi-user.target